# Tools/Cass/assets/openrc/cass
#!/sbin/openrc-run
// todo
command=/path/to/command
